# dot_config/systemd/user/ambxst-theme-propagate.service
[Unit]
Description=Propagate Ambxst theme preference to app toolkits

[Service]
Type=oneshot
ExecStart=%h/.local/bin/ambxst-propagate-theme

[Install]
WantedBy=default.target
