# dot_config/systemd/user/ambxst-theme-propagate.path
[Unit]
Description=Watch Ambxst theme preference

[Path]
PathChanged=%h/.config/ambxst/config/theme.json
PathChanged=%h/.config/ambxst/config
Unit=ambxst-theme-propagate.service

[Install]
WantedBy=default.target
